# roles/vault/templates/vault.service.j2
[Unit]
Description=vault service
Requires=network-online.target
After=network-online.target

[Service]
EnvironmentFile=-/etc/default/vault
Environment=GOMAXPROCS=2
Restart=on-failure
ExecStart=/usr/local/sbin/vault server $OPTIONS -config=/etc/vault.d/config.hcl
ExecReload=/bin/kill -HUP $MAINPID
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
